# deployments/quadlets/agent-os-api.container
[Unit]
Description=Agent OS Go API Backend
After=agent-os-db.service
Wants=agent-os-db.service

[Container]
Image=agent-os-api:latest
Network=agent-os.network
ContainerName=agent-os-api

Volume=/mnt/truenas/obsidian:/data/obsidian
Volume=/mnt/truenas/repos:/data/repos:ro
Volume=agent-os-artifacts:/data/artifacts
Volume=/home/tim/.hermes/skills:/data/hermes-skills:ro,z

EnvironmentFile=/home/tim/.config/containers/systemd/agent-os-api.env

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
